# perfume-shop.service
[Unit]
Description=Perfume Shop API Service
Documentation=https://github.com/yourusername/perfume-shop
After=syslog.target network-online.target remote-fs.target nss-lookup.target
Wants=network-online.target

[Service]
Type=simple
User=perfume
Group=perfume
WorkingDirectory=/opt/perfume-shop

# Environment file
EnvironmentFile=/opt/perfume-shop/.env.production

# Java executable path
ExecStart=/usr/bin/java \
  -Dspring.profiles.active=prod \
  -Dspring.config.location=classpath:application.yml,classpath:application-prod.yml \
  -Djava.io.tmpdir=/var/tmp/perfume-shop \
  ${JAVA_OPTS} \
  -jar perfume-shop-1.0.0.jar

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Security
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=true
ReadWritePaths=/var/log/perfume-shop

# Timeout settings
TimeoutStartSec=30
TimeoutStopSec=10

# Process management
StandardOutput=journal
StandardError=journal
SyslogIdentifier=perfume-shop

[Install]
WantedBy=multi-user.target
